# src/modules/mainsail/filesystem/root/etc/NetworkManager/dispatcher.d/99-r3dtospi5-nginx-root
#!/bin/sh
# Re-evaluate nginx "/" routing when a wired interface changes state.
IFACE="$1"
ACTION="$2"

case "${IFACE}" in
end0 | eth0) ;;
*) exit 0 ;;
esac

/usr/local/bin/r3dtospi5-update-nginx-root.sh --reload
exit 0
