# src/modules/mainsail/filesystem/root/etc/systemd/system/nginx.service.d/r3dtospi5-root.conf
[Service]
# Ensure /var/lib/r3dtospi5/nginx-mainsail-root.conf exists before nginx parses config.
ExecStartPre=/usr/local/bin/r3dtospi5-update-nginx-root.sh
